# ws/supervisord.conf
[program:ws]
command = ruby /home/dotcloud/current/server.rb
